# Installer/debian/ESET_backend.service
[Unit]
Description=ESET Backend Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/eset/client/backend/backend
Restart=always

[Install]
WantedBy=multi-user.target
